# service/firstrun/run
#!/bin/bash

ansible-playbook -i /etc/ansible/hosts /etc/ansible/itop.yml -c local

sv d firstrun

mv /etc/service/firstrun /root/
